# standalone/laser-keyboard.service
[Unit]
Description=laser-keyboard standalone (MIDI -> ArtNet)
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
# User and paths match the appliance install; adjust if you deploy elsewhere.
# A User= that doesn't exist on the host fails at boot with status=217/USER.
User=lichtfetisch
WorkingDirectory=/home/lichtfetisch/laser-keyboard/standalone
ExecStart=/home/lichtfetisch/laser-keyboard/standalone/.venv/bin/python -m laserkbd
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
